# deploy/crontab
TZ=Europe/Warsaw

# Wysyłanie kampanii emailowych
# Co 1 minute:
* * * * * cd /opt/mega-fun && /opt/mega-fun/venv/bin/python scripts/send_campaign.py >> /var/log/mega-fun-cron.log 2>&1

# Co 30 sekund (odkomentuj obie linie poniżej i zakomentuj powyższą):
# * * * * * cd /opt/mega-fun && /opt/mega-fun/venv/bin/python scripts/send_campaign.py >> /var/log/mega-fun-cron.log 2>&1
# * * * * * sleep 30 && cd /opt/mega-fun && /opt/mega-fun/venv/bin/python scripts/send_campaign.py >> /var/log/mega-fun-cron.log 2>&1

# Reset dziennego licznika emaili o północy
59 23 * * * cd /opt/mega-fun && /opt/mega-fun/venv/bin/python scripts/reset_daily_sent.py >> /var/log/mega-fun-cron.log 2>&1

# Backup bazy danych co noc o 3:00
0 3 * * * bash /opt/mega-fun/scripts/backup_db.sh